# miner.crontab
* * * * * bash /home/miner/reboot-on-fail.sh
