# miner.crontab
# M  H  DOM MO DOW CMD
*  *  *   *  *   bash /home/miner/reboot-on-fail.sh
12 12 *   *  *   sudo logrotate logrotate.conf
